# upstart/fluid/fluid_server.conf
description     "Fluid server"
author          "AMS"

start on (local-filesystems and net-device-up IFACE!=lo)

console log

exec /home/cloudlet/fluid-server/fluid_server
