# upstart/pycloud-api.conf
description     "Pycloud API"
author          "AMS"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid pycloud
setgid pycloud

console log

exec /usr/bin/pycloud-api
