# upstart/pycloud-manager.conf
description     "Pycloud Cloudlet Manager"
author          "AMS"

start on (local-filesystems and net-device-up IFACE!=lo)

console log

exec /usr/bin/pycloud-manager
